# src/etc/init.d/ixautomation-nat
#!/sbin/openrc-run
  
name="ixautomation-nat"
description="iXsystems automation framework"

depend()
{
  provide ixautomation-nat
  need dnsmasq
}

start()
{
  kldstat | grep -q ipfw_nat
  if [ $? -eq 1 ] ; then
    kldload ipfw_nat >/dev/null 2>/dev/null
    service ipfw restart >/dev/null 2>/dev/null
  fi
  sysctl net.inet.ip.forwarding | grep -q 1
  if [ $? -eq 1 ] ; then
    sysctl net.inet.ip.forwarding=1 >/dev/null 2>/dev/null
  fi
  ixautomation_bridge="ixautomation"
  if ! ifconfig ${ixautomation_bridge} >/dev/null 2>/dev/null ; then
    bridge=$(ifconfig bridge create)
    ifconfig ${bridge} name ${ixautomation_bridge} >/dev/null
    ifconfig ${ixautomation_bridge} 10.0.9.1 netmask 255.255.255.0 >/dev/null 2>/dev/null
  fi
  sysctl net.link.bridge.ipfw | grep -q 1
  if [ $? -eq 1 ] ; then
    sysctl net.link.bridge.ipfw=1 >/dev/null 2>/dev/null
  fi
  einfo "Created ixautomation bridge"
  inet=`netstat -nr | grep default | awk '{ print $4 }'`
  ipfw nat 9999 config if ${inet} >/dev/null 2>/dev/null
  ipfw add 9001 nat 9999 ip from 10.0.9.0/24 to any out via ${inet} >/dev/null 2>/dev/null
  ipfw add 9002 nat 9999 ip from any to any in via ${inet} >/dev/null 2>/dev/null
  ipfw add 9003 allow all from any to any >/dev/null 2>/dev/null
  einfo "NAT has been enabled for ixautomation bridge"
  vm init >/dev/null 2>/dev/null
}
